# bin/das_cacheserver
#!/usr/bin/env bash
#
# dbs_discovery This script runs CMS DBS Data Discovery service
#
# chkconfig: 345 05 95

if [ -z ${WMCORE_ROOT} ]; then
   echo "The WMCORE_ROOT environment is not set"
   exit 1
fi
export WTBASE=$WMCORE_ROOT/src
if [ -z $DAS_ROOT ]; then
   echo "The DAS_ROOT environment is not set"
   exit 1
fi

RETVAL=$?

cmd="$WMCORE_ROOT/src/python/WMCore/WebTools/Root.py --ini=$DAS_ROOT/src/python/DAS/web/das_cacheconfig.py"
pid=`ps auxw | grep das_cacheconfig | grep -v grep | awk 'BEGIN{ORS=" "} {print $2}'`

case "$1" in
 stdout)
        echo $"Checking for existing DAS service..."
        if [ ! -z "${pid}" ]; then
          kill -9 ${pid}
        fi
        echo $"Restart DAS service..."
        ${cmd}
        ;;
 restart)
        echo $"Checking for existing DAS service..."
        if [ ! -z "${pid}" ]; then
          kill -9 ${pid}
        fi
        echo $"Restart DAS service..."
        nohup ${cmd} 2>&1 1>& /dev/null < /dev/null &
        ;;
 start)
        if [ ! -z "${pid}" ]; then
          kill -9 ${pid}
        fi
        nohup ${cmd} 2>&1 1>& /dev/null < /dev/null &
        ;;
 status)
        if [ ! -z "${pid}" ]; then
          echo $"DAS is running, pid=${pid}"
          ps -f -wwww -p ${pid}
          exit 0
        fi
        echo $"DAS is stopped"
        exit 3
        ;;
 stop)
        if [ ! -z "${pid}" ]; then
          kill -9 ${pid}
        fi
        ;;
 *)
        echo $"Usage: $0 {start|stop|status|restart|stdout}"
        exit 1
        ;;
esac

exit $RETVAL

